# gemini-proxy/gemini-proxy.service
[Unit]
Description=Gemini Proxy - OpenAI-compatible proxy for Google Gemini
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/gemini-proxy
ExecStart=/usr/bin/node /opt/gemini-proxy/src/index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=8787

[Install]
WantedBy=multi-user.target
